# systemd/anding-cctv.service
[Unit]
Description=Anding CCTV Streaming Service
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/anding/anding-cctv-relay

# 시작 전 이미지 업데이트 (조용히)
ExecStartPre=/usr/bin/docker compose pull --quiet

# 자동 업데이트(watchtower) 포함해서 시작
ExecStart=/usr/bin/docker compose --profile auto-update up -d

# Graceful shutdown: 컨테이너 안전하게 종료
ExecStop=/usr/bin/docker compose --profile auto-update down --timeout 30

# 재시작 (설정 변경 후)
ExecReload=/usr/bin/docker compose --profile auto-update restart

User=anding
Group=docker
TimeoutStartSec=180
TimeoutStopSec=60

# 환경변수 파일 로드
EnvironmentFile=-/home/anding/anding-cctv-relay/.env

# 재시작 정책
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
